# infra/oracle-api.service
[Unit]
Description=TruthMachine Oracle API
Documentation=https://github.com/komapc/retro/blob/main/docs/ORACLE_API.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/truthmachine/api

# Gunicorn supervises uvicorn workers. The master process owns the listening
# socket on :8001 and keeps it open across SIGHUP, so `systemctl reload` swaps
# in new code with no dropped requests.
#
# Flags:
#   -k uvicorn.workers.UvicornWorker  ASGI-compatible worker class
#   --workers 2                       one serves while the other reloads
#   --graceful-timeout 30             give in-flight forecasts time to finish
#                                     (LLM chains can run 10-20s)
#   --timeout 120                     matches nginx proxy_read_timeout
#   --keep-alive 5                    matches nginx keepalive
#   (no --preload)                    preload forks from master and breaks
#                                     per-worker code reload on SIGHUP
ExecStart=/home/ubuntu/.local/bin/uv run gunicorn forecast_api.main:app \
    -k uvicorn.workers.UvicornWorker \
    --workers 2 \
    --bind 127.0.0.1:8001 \
    --graceful-timeout 30 \
    --timeout 120 \
    --keep-alive 5 \
    --access-logfile - \
    --error-logfile - \
    --log-level info

# Zero-downtime deploy: `sudo systemctl reload oracle-api` after `uv sync`.
# Gunicorn's master handles SIGHUP by spawning new workers with fresh imports
# and draining the old ones gracefully; the socket is never closed.
ExecReload=/bin/kill -HUP $MAINPID

# Full restart (e.g. on dep-graph changes that require new master) still works
# via `systemctl restart`, but incurs the usual 2-5s 502 window.
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/truthmachine/.env
StandardOutput=append:/home/ubuntu/truthmachine/oracle_log.txt
StandardError=append:/home/ubuntu/truthmachine/oracle_log.txt

[Install]
WantedBy=multi-user.target
