# infra/oracle-api.service
[Unit]
Description=TruthMachine Oracle API
Documentation=https://github.com/komapc/retro/blob/main/docs/ORACLE_API.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/truthmachine/api
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn forecast_api.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1 \
    --log-level info
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/truthmachine/.env
StandardOutput=append:/home/ubuntu/truthmachine/oracle_log.txt
StandardError=append:/home/ubuntu/truthmachine/oracle_log.txt

[Install]
WantedBy=multi-user.target
