# infra/truthmachine-poc.service
[Unit]
Description=TruthMachine PoC Pipeline Loop (Beat Polymarket)
After=network-online.target truthmachine.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/truthmachine
ExecStart=/bin/bash /home/ubuntu/truthmachine/infra/ec2_run_poc.sh
Restart=always
RestartSec=60
StandardOutput=append:/home/ubuntu/truthmachine/poc_log.txt
StandardError=append:/home/ubuntu/truthmachine/poc_log.txt

[Install]
WantedBy=multi-user.target
